# temp.service
[Unit]
Description=Log temperature to statsd

[Service]
Type=simple
ExecStart=/usr/local/bin/temp.py

[Install]
WantedBy=multi-user.target

